# mediamtx/rootfs/etc/services.d/mediamtx/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
bashio::log.warning "MediaMTX exited (exit code: $1)."
