# mediamtx/rootfs/etc/services.d/mediamtx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

readonly OPTIONS=/data/options.json
readonly CONFIG=/data/mediamtx.yml

# ── helpers ───────────────────────────────────────────────────────────────────

opt() { jq -r "$1 // empty" "${OPTIONS}"; }

bool() {
    local v
    v=$(jq -r "$1" "${OPTIONS}")
    [[ "${v}" == "true" ]] && echo "true" || echo "false"
}

# ── generate mediamtx.yml ─────────────────────────────────────────────────────

bashio::log.info "Generating MediaMTX configuration..."

cat > "${CONFIG}" <<YAML
# Generated by the MediaMTX Home Assistant add-on.
# Changes are overwritten on every restart – edit via the add-on Configuration tab.

###############################################
# General

logLevel: $(opt '.log_level')
logDestinations: [stdout]

###############################################
# Protocols

rtsp: $(bool '.rtsp')
rtspAddress: :8554

rtmp: $(bool '.rtmp')
rtmpAddress: :1935

hls: $(bool '.hls')
hlsAddress: :8888
hlsAlwaysRemux: $(bool '.hls_always_remux')

webrtc: $(bool '.webrtc')
webrtcAddress: :8889
webrtcLocalUDPAddress: :8189

srt: $(bool '.srt')
srtAddress: :8890

###############################################
# Control API

api: $(bool '.api')
apiAddress: :9997

YAML

# ── authentication ────────────────────────────────────────────────────────────

AUTH_ENABLED=$(jq -r '.auth_enabled' "${OPTIONS}")

if [[ "${AUTH_ENABLED}" == "true" ]]; then
    PUBLISH_USER=$(opt '.publish_user')
    PUBLISH_PASS=$(opt '.publish_pass')
    READ_USER=$(opt '.read_user')
    READ_PASS=$(opt '.read_pass')

    # Fall back to publish credentials for reading when no separate read user is set
    if [[ -z "${READ_USER}" ]]; then
        READ_USER="${PUBLISH_USER}"
        READ_PASS="${PUBLISH_PASS}"
    fi

    bashio::log.info "Authentication enabled (publish user: ${PUBLISH_USER})"

    cat >> "${CONFIG}" <<YAML
###############################################
# Authentication

authMethod: internal
authInternalUsers:
  - user: ${PUBLISH_USER}
    pass: ${PUBLISH_PASS}
    ips: []
    permissions:
      - action: publish
        path:
  - user: ${READ_USER}
    pass: ${READ_PASS}
    ips: []
    permissions:
      - action: read
        path:
      - action: playback
        path:
  # Allow API / metrics / pprof from localhost without credentials
  - user: any
    pass:
    ips: ['127.0.0.1', '::1']
    permissions:
      - action: api
      - action: metrics
      - action: pprof

YAML
else
    bashio::log.info "Authentication disabled (all clients accepted)"

    cat >> "${CONFIG}" <<'YAML'
###############################################
# Authentication

authMethod: internal
authInternalUsers:
  - user: any
    pass:
    ips: []
    permissions:
      - action: publish
        path:
      - action: read
        path:
      - action: playback
        path:
  - user: any
    pass:
    ips: ['127.0.0.1', '::1']
    permissions:
      - action: api
      - action: metrics
      - action: pprof

YAML
fi

# ── path defaults (recording) ─────────────────────────────────────────────────

RECORD=$(bool '.record')
RECORD_PATH=$(opt '.record_path')
RECORD_FORMAT=$(opt '.record_format')
RECORD_DELETE_AFTER=$(opt '.record_delete_after')

cat >> "${CONFIG}" <<YAML
###############################################
# Path defaults

pathDefaults:
  record: ${RECORD}
  recordPath: ${RECORD_PATH}
  recordFormat: ${RECORD_FORMAT}
  recordDeleteAfter: ${RECORD_DELETE_AFTER}

YAML

# ── custom config override (inserted before paths so paths: wins) ─────────────

CUSTOM=$(opt '.custom_config')
if [[ -n "${CUSTOM}" ]]; then
    bashio::log.info "Appending custom configuration block..."
    printf '\n###############################################\n# Custom configuration\n\n' >> "${CONFIG}"
    printf '%s\n' "${CUSTOM}" >> "${CONFIG}"
fi

# ── paths ─────────────────────────────────────────────────────────────────────

printf '\n###############################################\n# Paths\n\npaths:\n' >> "${CONFIG}"

PATHS_COUNT=$(jq '.paths | length' "${OPTIONS}")

if [[ "${PATHS_COUNT}" -gt 0 ]]; then
    while IFS= read -r path_json; do
        NAME=$(echo "${path_json}" | jq -r '.name')
        SOURCE=$(echo "${path_json}" | jq -r '.source // empty')
        SOURCE_ON_DEMAND=$(echo "${path_json}" | jq -r '.source_on_demand')
        PATH_RECORD=$(echo "${path_json}" | jq -r '.record')

        printf '  %s:\n' "${NAME}" >> "${CONFIG}"

        if [[ -n "${SOURCE}" ]]; then
            printf '    source: %s\n' "${SOURCE}" >> "${CONFIG}"
            printf '    sourceOnDemand: %s\n' "${SOURCE_ON_DEMAND}" >> "${CONFIG}"
        fi

        printf '    record: %s\n' "${PATH_RECORD}" >> "${CONFIG}"
    done < <(jq -c '.paths[]' "${OPTIONS}")
else
    # Default: accept any publisher on any path
    printf '  all_others:\n' >> "${CONFIG}"
fi

bashio::log.info "Configuration written to ${CONFIG}"
bashio::log.info "Starting MediaMTX..."

exec /usr/bin/mediamtx "${CONFIG}"
